# syncthing/rootfs/etc/init.d/syncthing
#!/bin/sh
### BEGIN INIT INFO
# Provides: syncthing-YJBeetle
# Required-Start: $local_fs $remote_fs
# Required-Stop: $local_fs $remote_fs
# Should-Start: $network
# Should-Stop: $network
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: daemonized version of syncthing.
# Description: Starts the syncthing daemon.
### END INIT INFO


DAEMON=/usr/bin/syncthing
USER=YJBeetle
PORT=8002


. /lib/lsb/init-functions


wait_process()
{
    while kill -0 "$1" 2> /dev/null; do
        sleep 0.5
    done
}


start()
{
    echo "Starting Syncthing"
    start-stop-daemon --start --quiet -b -o -c $USER -u $USER --exec $DAEMON --umask 0002 -- -gui-address=0.0.0.0:${PORT}
}


stop()
{
    start-stop-daemon --stop -o -c $USER -K -u $USER -x $DAEMON
    wait_process $(pidof $DAEMON)
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart|reload|force-reload)
      stop
      start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
    exit 1
    ;;
esac

exit 0
